# crontab/channel.sh
#!/bin/sh
cd /usr/share/nginx/html/weixin/
/usr/bin/php /usr/share/nginx/html/weixin/think channel
